# crontab
0 9,18 * * * sh -x /root/system_info.sh > /tmp/run.log
